# usr/lib/opkg/info/dms.list
/lib/isp_addr/isp_database.bin
/conf/model.xml
/conf/template/upnp.json
/bin/dms
/conf/template/cloud_config.json
/bin/dmclient
/conf/template/mwan_load_balance.json
/conf/template/wds.json
/etc/init.d/dms
/conf/template/dhcpd.json
/conf/template/wan_port_detect.json
/etc/mcbDesc.conf
/etc/plugin_config.json
/conf/template/plugin_config.json
/conf/template/uhttpd.json
/conf/template/wireless.json
/conf/template/cloud_status.json
/conf/template/network.json
/conf/template/lte_manager.json
/conf/template/firewall.json
/conf/cfgSyncBlackList.xml
/conf/template/iptv_manage.json
/etc/plugin_conflict.json
/conf/template/protocol.json
/conf/template/data_template_test.json
/conf/template/system.json
/lib/isp_addr/version
/etc/mcbJsonDesc.conf
/bin/smp.sh
