# usr/lib/opkg/info/kmod-qca-wifi-flash_16mb-profile.list
/ini/QCN9000.ini
/lib/wifi/qcacommands_ol_radio.xml
/ini/internal/QCA5018_i.ini
/ini/QCA5018.ini
/usr/sbin/thermaltool
/lib/firmware/IPQ5018/firmware_rdp_feature.ini
/lib/firmware/qcn9000/bdwlan.ba0
/lib/firmware/qcn9000/firmware_rdp_feature_512P.ini
/usr/lib/libwlanfw_wmi.so
/etc/init.d/wireless_start
/lib/modules/4.4.60/ath_pktlog.ko
/etc/config/qcacfg80211
/lib/firmware/qcn9000/fw_version.txt
/usr/sbin/apstats
/lib/firmware/IPQ5018/m3_fw.b00
/lib/firmware/IPQ5018/m3_fw.b01
/lib/firmware/IPQ5018/m3_fw.b02
/lib/modules/4.4.60/qdf.ko
/lib/wifi/qca_ftm_init.sh
/usr/sbin/spectraltool
/lib/wifi/tools_config
/lib/modules/4.4.60/qca_ol.ko
/ini/internal/global_i.ini
/lib/firmware/IPQ5018/q6_fw.mdt
/lib/firmware/qcn9000/qdss_trace_config.bin
/usr/sbin/cfg80211tool
/lib/firmware/IPQ5018/m3_fw.flist
/usr/lib/libqca_tools.so
/lib/firmware/qcn9000/firmware_rdp_feature.ini
/usr/sbin/radartool
/lib/modules/4.4.60/wifi_3_0.ko
/lib/wifi/qcawifi_countrycode.txt
/etc/init.d/boot-ftm
/usr/sbin/wifistats
/lib/firmware/qcn9000/amss.bin
/lib/firmware/IPQ5018/q6_fw.flist
/usr/sbin/wifitool
/ini/internal/QCN9000_i.ini
/ini/QCN9000_TP.ini
/lib/firmware/IPQ5018/bdwlan.b24
/lib/modules/4.4.60/asf.ko
/lib/firmware/qcn9000/m3.bin
/lib/wifi/qcawificfg80211.sh
/lib/firmware/IPQ5018/m3_fw.mdt
/usr/lib/libwlanfw_htt.so
/lib/wifi/qca-wifi-modules
/etc/init.d/diag_socket_app
/usr/sbin/wlanconfig
/lib/wifi/qcacommands_vap.xml
/lib/firmware/qcn9000/Data.msc
/lib/firmware/IPQ5018/qdss_trace_config.bin
/lib/firmware/IPQ5018/fw_version.txt
/ini/global.ini
/lib/firmware/IPQ5018/firmware_rdp_feature_512P.ini
/ini/QCA5018_TP.ini
/lib/modules/4.4.60/umac.ko
/lib/modules/4.4.60/qca_spectral.ko
/lib/firmware/IPQ5018/q6_fw.b00
/lib/firmware/IPQ5018/q6_fw.b01
/lib/firmware/IPQ5018/q6_fw.b02
/lib/firmware/IPQ5018/q6_fw.b03
/lib/firmware/IPQ5018/q6_fw.b04
/lib/firmware/IPQ5018/q6_fw.b05
/lib/firmware/IPQ5018/q6_fw.b07
/lib/firmware/IPQ5018/q6_fw.b08
/lib/firmware/IPQ5018/q6_fw.b09
/usr/sbin/pktlogconf
/lib/firmware/IPQ5018/q6_fw.b10
/lib/firmware/IPQ5018/q6_fw.b11
/lib/firmware/IPQ5018/q6_fw.b13
/lib/firmware/IPQ5018/q6_fw.b14
/lib/modules/4.4.60/mem_manager.ko
/lib/wifi/qcawifi.sh
