# usr/lib/opkg/info/qca-mcs-apps.list
/usr/sbin/mcsd
/etc/init.d/mcsd
/usr/sbin/mcsctl
